# systemd/user/battery-checker.service
[Unit]
Description=Check battery level

[Service]
Type=oneshot
ExecStart=%h/bin/battery-checker.sh
